# templates/workspace_services/guacamole/guacamole-server/docker/services/tomcat/run
#!/usr/bin/with-contenv sh
echo >&2 "starting tomcat"

# Only enable the Application Insights agent when a connection string is supplied.
if [ -n "${APPLICATIONINSIGHTS_CONNECTION_STRING:-}" ]; then
	export CATALINA_OPTS="${CATALINA_OPTS:-} -javaagent:/tmp/applicationinsights-agent.jar"
else
	echo >&2 "Application Insights disabled (missing APPLICATIONINSIGHTS_CONNECTION_STRING)"
fi
exec /usr/share/tomcat9/bin/catalina.sh run
