# templates/workspace_services/guacamole/guacamole-server/docker/services/guacd/run
#!/usr/bin/with-contenv sh
echo >&2 "starting guacd"
exec /opt/guacamole/sbin/guacd -f -b 0.0.0.0  -L $GUACD_LOG_LEVEL  -l 4822
